# borclone/systemd/borclone-back@.service
[Unit]
Description=Run backup part of borclone

[Service]
User=tom
ExecStart=/usr/local/bin/borclone back -e -c /etc/borclone.d/%i.conf
