# borclone/systemd/borclone-back@.service
[Unit]
Description=Run backup part of borclone

[Service]
User=tom
ExecStart=/usr/local/bin/borclone back -c /etc/borclone_%i.conf
