# borclone/systemd/borclone-back@.service
[Unit]
Description=Run backup part of borclone (%i)

[Service]
User=tom
ExecStart=/usr/local/bin/borclone back -c /etc/borclone.d/%i.conf
